# template/log-pp/log-pp.ini
[program:log-pp]
command=<:install_root:>/log-pp/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=5
environment=JAVA_HOME=<:java_home:>
