# template/permission/permission.ini
[program:permission]
command=<:install_root:>/permission/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
