# template/chartservice/chartservice.ini
[program:chart-service]
command=<:install_root:>/chartservice/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=20
