# template/elasticsearch/elasticsearch.ini
[program:elasticsearch]
command=<:install_root:>/elasticsearch/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startsecs=15
startretries=30
