# template/oneagent/oneagent.ini
[program:oneagent]
command=<:install_root:>/oneagent/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=10
