# template/startup/startup.ini
[program:start_services]
command=/bin/bash -c "<:install_root:>/startup/start_services.sh >> <:log_root:>/start_services.log 2>&1"
stopasgroup=true
autostart=true
autorestart=false
startretries=1000
startsecs=15
