# template/metric-proxy/metric-proxy.ini
[program:metric-proxy]
command=<:install_root:>/metric-proxy/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=8
