# template/cloudwiz-user/cloudwiz-user.ini
[program:cloudwiz-user]
environment=NODE_ENV="production"
command=<:install_root:>/node/bin/node <:install_root:>/cloudwiz-user/server.js
autostart=false
autorestart=true
startretries=1000
startsecs=15

