# template/umanager/umanager.ini
[program:umanager]
command=<:install_root:>/umanager/bin/um_daemon.py
stopasgroup=true
autostart=false
autorestart=true
startsecs=5
startretries=1000
