# template/opentsdb/opentsdb.ini
[program:opentsdb]
command=/bin/bash -c "<:install_root:>/opentsdb/bin/start.sh"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
