# template/python-daemon/python-daemon.ini
[program:python-daemon]
command=<:install_root:>/python-daemon/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=20
