# template/webfront/webfront.ini
[program:webfront]
command=<:install_root:>/webfront/bin/grafana-server -homepath=<:install_root:>/webfront
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
