# template/alertd/alertd.ini
[program:daemon]
command=<:install_root:>/alertd/bin/start.sh
stopasgroup=true
autostart=false
autorestart=unexpected
startretries=1000
startsecs=20
