# template/mysql/mysql.ini
[program:mysql-master]
command=<:install_root:>/mysql/bin/start.sh master
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15

[program:mysql-slave]
command=<:install_root:>/mysql/bin/start.sh slave
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
