# template/hadoop/hadoop.ini
[program:datanode]
command=/bin/bash -c "<:install_root:>/hadoop/bin/hdfs datanode"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HADOOP_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,JAVA_LIBRARY_PATH=<:install_root:>/hadoop/lib/native

[program:namenode]
command=/bin/bash -c "<:install_root:>/hadoop/bin/hdfs namenode"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HADOOP_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,JAVA_LIBRARY_PATH=<:install_root:>/hadoop/lib/native

[program:secondarynamenode]
command=/bin/bash -c "<:install_root:>/hadoop/bin/hdfs secondarynamenode"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HADOOP_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,JAVA_LIBRARY_PATH=<:install_root:>/hadoop/lib/native
