# template/hbase/hbase.ini
[program:master]
command=/bin/bash -c "<:install_root:>/hbase/bin/hbase-daemon.sh foreground_start master"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HBASE_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,HADOOP_HOME=<:install_root:>/hadoop

[program:backupmaster]
command=/bin/bash -c "<:install_root:>/hbase/bin/hbase-daemon.sh foreground_start master --backup"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HBASE_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,HADOOP_HOME=<:install_root:>/hadoop

[program:regionserver]
command=/bin/bash -c "<:install_root:>/hbase/bin/hbase-daemon.sh foreground_start regionserver"
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HBASE_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,HADOOP_HOME=<:install_root:>/hadoop

[program:hbase-thrift]
command=<:install_root:>/hbase/bin/hbase-daemon.sh foreground_start thrift
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=15
environment=JAVA_HOME=<:java_home:>,HBASE_CLASSPATH=<:install_root:>/hadoop/lib/hadoop-lzo-0.4.21-SNAPSHOT.jar,HADOOP_HOME=<:install_root:>/hadoop
