# template/jmxtrans/jmxtrans.ini
[program:jmxtrans]
command=<:install_root:>/jmxtrans/bin/start.sh
stopasgroup=true
autostart=false
autorestart=true
startretries=1000
startsecs=20
